# pistelink/deploy/systemd/pistelink.service
# PisteLink backend service for the product install root.
[Unit]
Description=PisteLink backend (bare systemd)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pistelink-backend
Group=pistelink
SupplementaryGroups=dialout audio video render
WorkingDirectory=/opt/pistelink
Environment=PISTELINK_CONFIG=/etc/pistelink/config.toml
Environment=PISTELINK_RUNTIME_DIR=/run/pistelink
Environment=PISTELINK_FRONTEND_DIR=/opt/pistelink/frontend/dist
Environment=PISTELINK_SOUND_DIR=/opt/pistelink/sound
RuntimeDirectory=pistelink
RuntimeDirectoryMode=0770
RuntimeDirectoryPreserve=yes
UMask=0007
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
ReadWritePaths=/var/lib/pistelink /run/pistelink
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
MemoryDenyWriteExecute=false
ExecStart=/opt/pistelink/run_pistelink_backend.sh
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
