# pistelink/deploy/systemd/pistelink.service
# PisteLink backend service for this Jetson workspace.
# Host setup:
#   sudo /home/thomas/fencing/pistelink/deploy/install.sh
#   sudo cp /home/thomas/fencing/pistelink/deploy/systemd/pistelink.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now pistelink.service
[Unit]
Description=PisteLink backend (bare systemd)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=thomas
Group=thomas
SupplementaryGroups=dialout audio video render
WorkingDirectory=/home/thomas/fencing/pistelink
Environment=PISTELINK_CONFIG=/etc/pistelink/config.toml
Environment=PISTELINK_RUNTIME_DIR=/run/pistelink
Environment=PISTELINK_FRONTEND_DIR=/home/thomas/fencing/pistelink/frontend/dist
Environment=PISTELINK_SOUND_DIR=/home/thomas/fencing/pistelink/sound
RuntimeDirectory=pistelink
RuntimeDirectoryMode=0770
RuntimeDirectoryPreserve=yes
UMask=0007
ExecStart=/home/thomas/fencing/pistelink/run_pistelink_backend.sh
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
