# pistelink_ai_pipeline/deploy/systemd/pistelink-ai.service
# PisteLink AI camera/judging service for the product install root.
[Unit]
Description=PisteLink AI camera and judging service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pistelink-ai
Group=pistelink
SupplementaryGroups=video audio render dialout
WorkingDirectory=/opt/pistelink/ai/service
Environment=PISTELINK_AI_SOCKET=/run/pistelink/ai.sock
Environment=PISTELINK_RUNTIME_DIR=/run/pistelink
Environment=PISTELINK_AI_SOCKET_GROUP=pistelink
Environment=PISTELINK_AI_SOCKET_MODE=0660
Environment=PISTELINK_MATCH_ROOT=/var/lib/pistelink/matches
Environment=PISTELINK_ANALYZER_ROOT=/opt/pistelink/ai/analyzer
Environment=PISTELINK_ANALYZER_MODEL_PATH=/opt/pistelink/ai/models/yolo26l-pose_fast_fp16_ultra.engine
Environment=PISTELINK_BLADE_MODEL_PATH=/opt/pistelink/ai/models/blade_touch_referee_model.joblib
Environment=PISTELINK_ANALYZER_FISHEYE_BACKEND=none
Environment=PISTELINK_ANALYZER_STARTUP_TIMEOUT=120
Environment=FENCING_CAMERA_WIDTH=1280
Environment=FENCING_CAMERA_HEIGHT=720
Environment=FENCING_CAMERA_FPS=30
Environment=PYTHONNOUSERSITE=1
RuntimeDirectory=pistelink
RuntimeDirectoryMode=0770
RuntimeDirectoryPreserve=yes
UMask=0007
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
ReadWritePaths=/var/lib/pistelink /run/pistelink
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
MemoryDenyWriteExecute=false
ExecStart=/opt/pistelink/ai/service/run_pistelink_ai_service.sh
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
